# docker-cleaner/files/99-docker-cleanup.cron.daily
#!/bin/sh

docker rm $(docker ps -aqf status=exited)
docker rmi $(docker images -qf dangling=true)
docker volume rm $(docker volume ls -qf dangling=true)

docker run --rm -e GRACE_PERIOD_SECONDS=86400 -e FORCE_IMAGE_REMOVAL=1 -v /var/run/docker.sock:/var/run/docker.sock spotify/docker-gc
